# infra/equipecoo.service
[Unit]
Description=Equipes Coo — FastAPI Backend
After=network.target

[Service]
User=dudox
WorkingDirectory=/home/dudox/equipecoo/backend
EnvironmentFile=/home/dudox/equipecoo/backend/.env
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
